# src/tools/secrets/otra_tpm.socket.dist
[Unit]
Description=OTRA TPM Daemon Socket

[Socket]
ListenStream=/run/otra/tpm_daemon.sock
SocketMode=0770
SocketUser=myuser
SocketGroup=www-data

[Install]
WantedBy=sockets.target
